# root/usr/lib/systemd/system/post-system-init.target
[Unit]
Description=NethServer is configured
Documentation=https://github.com/nethesis/nethserver-base
After=nethserver-system-init.service